# References/lvmonluks/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin rose --noclear %I $TERM
Type=simple
